# systemd/wxconsumer.service
[Unit]
Description=wxconsumer
After=network.target
 
[Service]
Type=forking
ExecStart=/app/wxpark/consumer.sh
 
[Install]
WantedBy=multi-user.target
