# systemd/wxgateway.service
[Unit]
Description=wxgateway
After=network.target
 
[Service]
Type=forking
ExecStart=/app/wxpark/gateway.sh
 
[Install]
WantedBy=multi-user.target
