# systemd/wxip.service
[Unit]
Description=wxip
After=network.target
 
[Service]
Type=forking
ExecStart=/app/park-script/ip.sh
 
[Install]
WantedBy=multi-user.target
