# systemd/wxjob.service
[Unit]
Description=wxjob
After=network.target
 
[Service]
Type=forking
ExecStart=/app/park-script/job.sh
 
[Install]
WantedBy=multi-user.target
