# systemd/wxprovider.service
[Unit]
Description=wxprovider
After=network.target
 
[Service]
Type=forking
ExecStart=/app/wxpark/provider.sh
 
[Install]
WantedBy=multi-user.target
